# attacker/supervisord.conf
[supervisord]
nodaemon=true

[program:keep_alive]
command=/keep_alive.sh
autostart=true
autorestart=true
